# docker/slurmctld/etc/supervisor/conf.d/sssd.conf
[program:sssd]
command=/usr/sbin/sssd -i --logger=files
pidfile=/run/sssd.pid
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
